# caddy/rootfs/etc/services.d/caddy/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Community Hass.io Add-ons: Caddy
# Runs Caddy
# ==============================================================================

exec caddy -conf '/share/addon-caddy/Caddyfile'